# profile.d/python.sh
if test -f $HOME/.pythonstartup.py; then
export PYTHONSTARTUP=$HOME/.pythonstartup.py
fi
