# profile.d/vim-ish.sh
alias vim='vim -X'

# use vi mode
set -o vi

export EDITOR=vim
